# ansible/templates/lustre-modules.service.j2
[Unit]
Description=Load Lustre Modules Safely
After=network-online.target
Wants=network-online.target
Conflicts=shutdown.target

[Service]
Type=oneshot
# Bring up LNET interfaces before loading lnet so ksocklnd binds to the
# correct static IPs. network-online.target fires on eth0/DHCP only;
# lustre-net0/lustre-net1 (eth1/eth2) may not be configured yet without this.
ExecStartPre=/usr/bin/nmcli connection up lustre-net0
{% if ip_lnet1 is defined %}
ExecStartPre=/usr/bin/nmcli connection up lustre-net1
{% endif %}
ExecStart=/sbin/modprobe lnet
ExecStart=/sbin/modprobe lustre
RemainAfterExit=yes
TimeoutSec=60

[Install]
WantedBy=multi-user.target
