# ansible/templates/lustre-modules.service.j2
[Unit]
Description=Load Lustre Modules Safely
After=network-online.target
Wants=network-online.target
Conflicts=shutdown.target

[Service]
Type=oneshot
# Bring up the LNET interface before loading lnet so ksocklnd binds to the
# correct static IP. network-online.target fires on eth0/DHCP only;
# lustre-net0 (eth1) may not be configured yet without this.
ExecStartPre=/usr/bin/nmcli connection up lustre-net0
ExecStart=/sbin/modprobe lnet
ExecStart=/sbin/modprobe lustre
RemainAfterExit=yes
TimeoutSec=60

[Install]
WantedBy=multi-user.target
